# infra/cluster-recall/cluster-recall.service
[Unit]
Description=Cluster Recall HTTP service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cluster-recall
ExecStart=/usr/bin/python3 /opt/cluster-recall/cluster_recall.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
