# infra/cluster-guardian/cluster-guardian.timer
[Unit]
Description=Cluster Guardian alle 60s

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=5s

[Install]
WantedBy=timers.target
