# infra/cluster-guardian/cluster-guardian.service
[Unit]
Description=Cluster Guardian L0 (lokaler Floor)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/cluster-guardian/cluster-guardian.sh
TimeoutStartSec=30
